# crontab/updateCode.sh
#!/bin/bash
###############################################################
# FileName: updateCode.sh
# Author: Chao Hong <teddy.hongchao@gmail.com>
# Program:
###############################################################
php /$yourPath/gitlab_project/bin/init.php > /dev/null
